# os/linux-2.6.21.5/root/usr/lib/ipkg/info/base-files-atheros-2.6.list
/
/var
/dev/
/bin/
/bin/uci
/bin/ipcalc.sh
/bin/login
/bin/firstboot
/lib/
/init
/etc/
/etc/protocols
/etc/inittab
/etc/sysctl.conf
/etc/mtab
/etc/profile
/etc/config/
/etc/config/system
/etc/config/network
/etc/banner
/etc/functions.sh
/etc/rc.d/
/etc/init.d/
/etc/init.d/done
/etc/init.d/umount
/etc/init.d/usb
/etc/init.d/sysctl
/etc/init.d/rcS
/etc/init.d/boot
/etc/preinit.arch
/etc/group
/etc/crontabs/
/etc/uci-defaults/
/etc/diag.sh
/etc/resolv.conf
/etc/rc.sh
/etc/preinit
/etc/hosts
/etc/shells
/etc/rc.common
/etc/passwd
/etc/ipkg.conf
/mnt/
/sbin/
/sbin/ifdown
/sbin/usb-storage
/sbin/hotplug-call
/sbin/mount_root
/sbin/wifi
/sbin/ifup
/www/
/usr/
/usr/bin/
/usr/lib/
/usr/lib/common.awk
/usr/share/
/usr/share/udhcpc/
/usr/share/udhcpc/default.script
/tmp/
/sys/
/proc/
