# os/linux-2.6.21.5/root/usr/lib/ipkg/info/busybox.list
/
/bin/
/bin/pwd
/bin/netstat
/bin/gzip
/bin/kill
/bin/vi
/bin/echo
/bin/mount
/bin/chmod
/bin/mknod
/bin/mkdir
/bin/df
/bin/ip
/bin/dd
/bin/ping
/bin/sh
/bin/ps
/bin/ash
/bin/umount
/bin/zcat
/bin/dmesg
/bin/ls
/bin/sync
/bin/rmdir
/bin/chown
/bin/sleep
/bin/mv
/bin/busybox
/bin/cat
/bin/fgrep
/bin/ln
/bin/cp
/bin/touch
/bin/rm
/bin/egrep
/bin/tar
/bin/grep
/bin/pidof
/bin/gunzip
/etc/
/etc/init.d/
/etc/init.d/cron
/etc/init.d/telnet
/etc/init.d/httpd
/sbin/
/sbin/rmmod
/sbin/reboot
/sbin/route
/sbin/udhcpc
/sbin/init
/sbin/halt
/sbin/lsmod
/sbin/insmod
/sbin/klogd
/sbin/sysctl
/sbin/poweroff
/sbin/ifconfig
/sbin/logread
/sbin/syslogd
/sbin/getty
/usr/
/usr/bin/
/usr/bin/awk
/usr/bin/find
/usr/bin/printf
/usr/bin/killall5
/usr/bin/unlzma
/usr/bin/arping
/usr/bin/logger
/usr/bin/uptime
/usr/bin/tftp
/usr/bin/[
/usr/bin/free
/usr/bin/du
/usr/bin/lzmacat
/usr/bin/test
/usr/bin/awx
/usr/bin/top
/usr/bin/bzcat
/usr/bin/[[
/usr/bin/ftpget
/usr/bin/bunzip2
/usr/bin/cut
/usr/bin/reset
/usr/bin/cksum
/usr/bin/killall
/usr/bin/ftpput
/usr/bin/which
/usr/sbin/
/usr/sbin/telnetd
/usr/sbin/inetd
/linuxrc
